# debian/bacula-director-sqlite3.prerm
#! /bin/sh

set -e


. /usr/share/debconf/confmodule
db_version 2.0


THISDB=sqlite3
if [ -f /usr/share/dbconfig-common/dpkg/prerm.$THISDB ]; then
	. /usr/share/dbconfig-common/dpkg/prerm.$THISDB
	dbc_go bacula-director-$THISDB $@
fi

if test "$1" = "remove" -o "$1" = "deconfigure"; then
  invoke-rc.d bacula-director stop || exit $?-
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
